# target/skeleton/etc/inittab
# Start the system
::sysinit:/etc/init.d/rc.S

# Start a terminal on the serial port's
ttyS0::respawn:/sbin/getty -L 115200 /dev/ttyS0 vt100

# What to do on graceful shutdown
::shutdown:/etc/init.d/rc.0

# What to do when restarting init
::restart:/sbin/init
