# target/overlay/etc/udev/rules.d/90-storage-hotplug.rules
ENV{ID_FS_USAGE}=="filesystem", RUN+="/sbin/hotplug.storage"
